# examples/exten_gunicorn.socket
[Unit]
Description=exten socket

[Socket]
ListenStream=/run/exten.sock

[Install]
WantedBy=sockets.target